# systemd/launch-week-monitor.timer
[Unit]
Description=BuyWhere launch week health monitor (once per 5 min)
After=network-online.target
Wants=network-online.target

[Timer]
OnCalendar=*:0/5:00
AccuracySec=1s
Persistent=true
Unit=launch-week-monitor.service

[Install]
WantedBy=timers.target